# src/etc/services-config/supervisor/supervisord.d/dnsmasq-wrapper.conf
[program:dnsmasq-wrapper]
priority = 100
command = /usr/sbin/dnsmasq-wrapper
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/dnsmasq.log
stdout_events_enabled = true